# home/Library/LaunchAgents/dev.kryota.morning-radar.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Weekday-morning radar (kryota-dev/dotfiles#257).
  Fires the morning-brief headless wrapper at 09:00 local time, Monday through
  Friday. Local time is assumed to be JST on this machine. launchd coalesces
  fires missed while asleep into one run on wake; days the Mac is powered off
  are skipped entirely. Registered/reloaded by
  run_onchange_after_30-register-launchd-agents.sh.tmpl.
  RunAtLoad is intentionally absent (defaults to false) so registration itself
  never triggers a billed run.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.kryota.morning-radar</string>
  <key>ProgramArguments</key>
  <array>
    <string>{{ .chezmoi.homeDir }}/.claude/morning-radar.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>2</integer>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>3</integer>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>4</integer>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>5</integer>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>
  <key>WorkingDirectory</key>
  <string>{{ .chezmoi.homeDir }}/dotfiles</string>
  <key>StandardOutPath</key>
  <string>{{ .chezmoi.homeDir }}/Library/Logs/dev.kryota.morning-radar.log</string>
  <key>StandardErrorPath</key>
  <string>{{ .chezmoi.homeDir }}/Library/Logs/dev.kryota.morning-radar.log</string>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
